# devel/.private/ign_migration_scripts/etc/catkin/profile.d/29.ign_migration_scripts.zsh
# Export Ignition environment variables

export IGN_LAUNCH_CONFIG_PATH=/home/admin-vras/aloam_ws/src/subt/ign_migration_scripts/launch
export IGN_GAZEBO_RESOURCE_PATH=/home/admin-vras/aloam_ws/src/subt/ign_migration_scripts/worlds:${IGN_GAZEBO_RESOURCE_PATH}
