# devel/.private/subt_ign/etc/catkin/profile.d/30.subt_ign.zsh
# Export Ignition environment variables

export IGN_LAUNCH_CONFIG_PATH=/home/admin-vras/aloam_ws/src/subt/subt_ign/launch
export IGN_LAUNCH_PLUGIN_PATH=/home/admin-vras/aloam_ws/devel/.private/subt_ign/lib
export IGN_GAZEBO_RESOURCE_PATH=/home/admin-vras/aloam_ws/src/subt/subt_ign/worlds
export SUBT_IMAGES_PATH=/home/admin-vras/aloam_ws/src/subt/subt_ign/images
export IGN_TRANSPORT_TOPIC_STATISTICS=1
